# scripts/com.rwill128.algorithm-visualization.youtube-upload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rwill128.algorithm-visualization.youtube-upload</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "/Users/richardwilliams/algorithm-visualization-uploader" || exit 1; mkdir -p logs; { echo "[$(/bin/date -u '+%Y-%m-%dT%H:%M:%SZ')] youtube upload queue run"; args=(scripts/youtube_upload_queue.py upload-next --limit "${UPLOAD_LIMIT:-1}"); if [[ "${DRY_RUN:-0}" == "1" ]]; then args+=(--dry-run); else args+=(--approved --allow-public); fi; .venv/bin/python "${args[@]}"; if [[ "${DRY_RUN:-0}" != "1" ]] &amp;&amp; ! /usr/bin/git diff --quiet -- data/youtube_upload_queue.json; then /usr/bin/git add data/youtube_upload_queue.json; /usr/bin/git commit -m "Record YouTube upload queue state"; /usr/bin/git push origin main; fi; } &gt;&gt; logs/youtube-upload-queue.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>11</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>13</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>17</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/richardwilliams/algorithm-visualization-uploader</string>
  <key>StandardOutPath</key>
  <string>/Users/richardwilliams/algorithm-visualization-uploader/logs/youtube-upload-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/richardwilliams/algorithm-visualization-uploader/logs/youtube-upload-launchd.err.log</string>
</dict>
</plist>
